# docker/etc/supervisord.d/processes.ini
[program:gunicorn]
command=/opt/lcogt-python36/bin/gunicorn --bind=0.0.0.0:8080 --workers=4 --access-logfile=- --error-logfile=- neox.wsgi:application
directory=/var/www/apps/neoexchange
user=nginx
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
