# config/processes.ini
[program:uwsgi]
command=/usr/bin/uwsgi --ini /etc/uwsgi.ini
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:celerymain]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command = celery -A neox.celeryconf worker -Q celery -n main_worker -f /dev/null --workdir=/var/www/apps/neoexchange/

[program:celeryperiodic]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command = celery -A neox.celeryconf worker -Q periodic -n periodic_worker --concurrency=1 -f /dev/null --workdir=/var/www/apps/neoexchange/

[program:celerybeat]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command = celery -A neox.celeryconf beat --workdir=/var/www/apps/neoexchange/

[program:mpcstatus]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
command = python /var/www/apps/neoexchange/neox/socket.py
